# development/net-analyzer/httpry/files/httpry.rc.1
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/ntp/files/ntpd.rc,v 1.24 2010/07/18 21:53:24 vapier Exp $

depend() {
	need net
}

start() {
	ebegin "Starting httpry"
	start-stop-daemon --start --quiet --exec /usr/sbin/httpry \
		--pidfile /var/run/httpry.pid \
		-- -d -q "${HTTPRY_OPTS}" >/dev/null 2>&1
	eend $? "Failed to start httpry"

}

stop() {
	ebegin "Stopping httpry"
	start-stop-daemon --stop --quiet --pidfile /var/run/httpry.pid
	eend $? "Failed to stop httpry"
}

